# validator-automation/systemd/meshchain-validators.service
[Unit]
Description=MeshChain local validator set (automation package)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=meshchain
Group=meshchain
WorkingDirectory=/opt/meshchain
Environment=BUILD_PROFILE=release
Environment=CHAIN_MODE=host
Environment=AUTO_BUILD=0
ExecStart=/opt/meshchain/validator-automation/mesh-validator --mode host --release start
ExecStop=/opt/meshchain/validator-automation/mesh-validator --mode host stop
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
