# deploy/meshchain-testnet.service
[Unit]
Description=MeshChain public testnet host (3 validators + faucet + scanner)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=meshchain
Group=meshchain
WorkingDirectory=/opt/meshchain
EnvironmentFile=-/etc/meshchain/testnet.env
ExecStart=/usr/local/bin/meshchain-testnet-start
ExecStop=/usr/local/bin/meshchain-testnet-stop
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
